# crontab
* * * * *   root   php /app/artisan schedule:run
* * * * *   root   echo "Running"
