# crontab
* * * * *   root   php /app/artisan schedule:run
